# deploy/collector/cc.kiki.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchDaemon for headless macOS — installs over SSH with sudo.

  IMPORTANT: code MUST live OUTSIDE ~/Documents — modern macOS TCC blocks
  daemons from user data directories ("Current directory does not exist" in
  the log) even when the daemon runs as that user. We use /opt/kiki-collector/.

  One-time setup on studio:
    sudo mkdir -p /opt/kiki-collector
    sudo chown clems:staff /opt/kiki-collector
    rsync -az ~/Documents/Projets/kiki-cockpit/deploy/collector/ /opt/kiki-collector/
    cd /opt/kiki-collector && ~/.local/bin/uv sync

    sudo cp /opt/kiki-collector/cc.kiki.collector.plist /Library/LaunchDaemons/
    sudo chown root:wheel /Library/LaunchDaemons/cc.kiki.collector.plist
    sudo chmod 644       /Library/LaunchDaemons/cc.kiki.collector.plist
    sudo launchctl bootout   system/cc.kiki.collector 2>/dev/null
    sudo launchctl bootstrap system /Library/LaunchDaemons/cc.kiki.collector.plist

  Verify:
    curl http://localhost:9150/healthz
    curl http://studio:9150/healthz       # from another tailnet member
    tail -f ~/Library/Logs/kiki-collector.log

  Reload after editing the plist:
    sudo cp /opt/kiki-collector/cc.kiki.collector.plist /Library/LaunchDaemons/
    sudo launchctl bootout   system/cc.kiki.collector
    sudo launchctl bootstrap system /Library/LaunchDaemons/cc.kiki.collector.plist
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cc.kiki.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/clems/.local/bin/uv</string>
        <string>run</string>
        <string>--directory</string>
        <string>/opt/kiki-collector</string>
        <string>uvicorn</string>
        <string>kiki_collector.main:app</string>
        <string>--host</string>
        <string>0.0.0.0</string>
        <string>--port</string>
        <string>9150</string>
    </array>

    <key>UserName</key>
    <string>clems</string>
    <key>GroupName</key>
    <string>staff</string>

    <key>WorkingDirectory</key>
    <string>/opt/kiki-collector</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>COLLECTOR_MACHINE_LABEL</key>
        <string>studio</string>
        <key>HOME</key>
        <string>/Users/clems</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/clems/Library/Logs/kiki-collector.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/clems/Library/Logs/kiki-collector.log</string>

    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
